# Library/LaunchAgents/com.mbezemer.sync-kubeconfig.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mbezemer.sync-kubeconfig</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{ .chezmoi.homeDir }}/bin/sync-kubeconfig-to-1password</string>
  </array>

  <key>WatchPaths</key>
  <array>
    <string>{{ .chezmoi.homeDir }}/.kube/config</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/sync-kubeconfig.log</string>

  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/sync-kubeconfig.log</string>
</dict>
</plist>
